# trudy-backend.service
[Unit]
Description=Trudy Backend API Service
After=network.target

[Service]
Type=notify
User=root
Group=root
WorkingDirectory=/opt/trudy-backend
Environment="PATH=/opt/trudy-backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/trudy-backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trudy-backend

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/trudy-backend /mnt/storage

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
